# aroma-diffuser.service
[Unit]
Description=Oil Diffuser Control System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/aroma-pi
ExecStart=/home/pi/aroma-pi/venv/bin/python /home/pi/aroma-pi/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

